# consul-monitor/check_consul-supervisord.conf
[program:check_consul]
command=/check_consul.sh
autostart=true
autorestart=true