# consul-monitor/consul-supervisord.conf
[program:consul]
command=/consul_start.sh
autostart=true
autorestart=true